# roles/lightbulb-ansiblered-deck/templates/ansible-php-deck.service
[Unit]
Description=PHP daemon to serve the Ansible deck

[Service]
ExecStart=/root/ansible-php-deck.sh

[Install]
WantedBy=multi-user.target

